# deploy/camilladsp.service
[Unit]
Description=CamillaDSP for Coldth
After=sound.target
Before=coldth.service

[Service]
Type=simple
User=livingroom
Group=livingroom
SupplementaryGroups=audio
WorkingDirectory=/home/livingroom
ExecStart=/usr/local/bin/camilladsp -s /home/livingroom/camilladsp/statefile.yml -w -g0 -o /home/livingroom/camilladsp/camilladsp.log -p 1234
Restart=always
RestartSec=2
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/livingroom/camilladsp

[Install]
WantedBy=multi-user.target
