# deploy/coldth.service
[Unit]
Description=Coldth web equalizer
Wants=network-online.target
After=network-online.target
Before=camilladsp.service

[Service]
Type=simple
User=coldth
Group=coldth
SupplementaryGroups=audio
WorkingDirectory=/opt/coldth
Environment=COLDTH_DATA_DIR=/var/lib/coldth
Environment=COLDTH_HOST=0.0.0.0
Environment=COLDTH_PORT=8080
Environment=COLDTH_CAPTURE_DEVICE=hw:Loopback,1,0
Environment=COLDTH_PLAYBACK_DEVICE=hw:Headphones,0
Environment=COLDTH_CAPTURE_FORMAT=S16LE
Environment=COLDTH_PLAYBACK_FORMAT=S16LE
ExecStart=/opt/coldth/.venv/bin/coldth
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# Keep home directories read-only so this unit also works when a user chooses
# a home-directory checkout instead of /opt/coldth.
ProtectHome=read-only
ReadWritePaths=/var/lib/coldth

[Install]
WantedBy=multi-user.target
